# templates/restic.sh.timer.j2
[Unit]
Description=Regulary restic backup

[Timer]
OnCalendar={{ restic_backup_timespan }}
[Install]
WantedBy=basic.target
