# files/restic.sh.service
[Unit]
Description=Restic Backup Job
[Service]
Type=oneshot
ExecStart=/usr/local/bin/restic.sh
KillMode=process

